# deploy/mhr-bot.service
# systemd サービス定義: モンハン傀異化素材検索 Discord Bot
#
# 使い方はリポジトリ直下の「ラズパイ移行ガイド.md」を参照。
# ※ User と パスの「pi」の部分は、自分のラズパイのユーザー名に合わせて書き換えてください。

[Unit]
Description=MHR Kairika Material Search Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/mhr-bot
ExecStart=/home/pi/mhr-bot/.venv/bin/python bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
